# user/systemd/user/waybar.service.d/sway-session.conf
[Unit]
PartOf=sway-session.target
Before=sway-session.target

[Install]
WantedBy=sway-session.target
